# systemd/hibp-checker.service
[Unit]
Description=HIBP Breach Checker
Documentation=https://github.com/greogory/hibp-checker
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%u
Group=%u

# IMPORTANT: Update this path to your hibp-project installation directory
# Example: WorkingDirectory=%h/hibp-project
WorkingDirectory=HIBP_PROJECT_DIR

# Environment
Environment="HIBP_API_KEY=%E"
EnvironmentFile=-%h/.config/hibp-checker/hibp-checker.env

# Execution
# Exit codes: 0=no breaches, 1=breaches found, 2=critical (password exposed)
# All are valid outcomes, not failures
SuccessExitStatus=1 2

# IMPORTANT: Update this path to match WorkingDirectory
ExecStart=HIBP_PROJECT_DIR/hibp_workflow.sh check

# Output (XDG-compliant location)
StandardOutput=append:%h/.local/share/hibp-checker/hibp-checker.log
StandardError=append:%h/.local/share/hibp-checker/hibp-checker.error.log

# Security
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/hibp-checker

[Install]
WantedBy=default.target
